# deployment/systemd/duckdns-update.service
[Unit]
Description=DuckDNS Updater (KeukaSensorProd)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/home/pi/KeukaSensorProd

# Ensure directories exist before running
ExecStartPre=/bin/mkdir -p /home/pi/KeukaSensorProd/configuration/services /home/pi/KeukaSensorProd/data/logs

# Environment allows overrides without editing the script
Environment=DATA_DIR=/home/pi/KeukaSensorProd/data
Environment=CONFIG_DIR=/home/pi/KeukaSensorProd/configuration/services
Environment=LOG_DIR=/home/pi/KeukaSensorProd/data/logs
Environment=DUCKDNS_CONFIG=/home/pi/KeukaSensorProd/configuration/services/duckdns.conf
Environment=DUCKDNS_TMP=/home/pi/KeukaSensorProd/configuration/services/duckdns.tmp
Environment=DUCKDNS_LOG=/home/pi/KeukaSensorProd/data/logs/duckdns_last.txt

ExecStart=/bin/bash /home/pi/KeukaSensorProd/deployment/scripts/duckdns_update.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
