# systemd/duckdns-update.service
[Unit]
Description=DuckDNS Update Script
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/home/pi/KeukaSensorProd
# Run with /bin/sh (dash) to match the script's shebang assumptions
ExecStart=/bin/sh /home/pi/KeukaSensorProd/scripts/duckdns_update.sh
# Optional overrides (e.g., DUCKDNS_CONF/LOG/IPV6)
EnvironmentFile=-/etc/default/duckdns
# Light hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
TasksMax=16
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target