# systemd/duckdns-update.service
[Unit]
Description=DuckDNS Update Script
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
# Optional env overrides (DUCKDNS_CONF/LOG/IPV6/CURL) in /etc/default/duckdns
EnvironmentFile=-/etc/default/duckdns
WorkingDirectory=/home/pi/KeukaSensorProd
# IMPORTANT: run with Bash (avoids sh option incompatibilities)
ExecStart=/usr/bin/env bash /home/pi/KeukaSensorProd/scripts/duckdns_update.sh
# Treat exit code 1 (DuckDNS KO) as SUCCESS so the unit doesn't "fail"
SuccessExitStatus=1

[Install]
WantedBy=multi-user.target
