# systemd/duckdns-update.service
[Unit]
Description=DuckDNS one-shot updater
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /home/pi/KeukaSensorProd/scripts/duckdns_update.sh
