# systemd/duckdns-update.timer
[Unit]
Description=Run DuckDNS updater periodically

[Timer]
# First run shortly after boot, then every 5 minutes
OnBootSec=2min
OnUnitActiveSec=5min
Unit=duckdns-update.service
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
