# systemd/duckdns-update.service
[Unit]
Description=DuckDNS Update Script
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/home/pi/KeukaSensorProd

# Explicit environment so it matches your successful manual run
Environment=DUCKDNS_CONF=/opt/keuka-sensor/duckdns.conf
Environment=DUCKDNS_LOG=/home/pi/KeukaSensorProd/logs/duckdns_last.txt
Environment=CURL=/usr/bin/curl
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/pi

# Run with Bash explicitly
ExecStart=/usr/bin/env bash /home/pi/KeukaSensorProd/scripts/duckdns_update.sh

# Treat remote KO (exit 1) as success so the unit doesn't look failed
SuccessExitStatus=1

# Send stdout/stderr to the journal for easier debugging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
