# systemd/keuka-sensor.service
[Unit]
Description=Keuka Sensor web app (Gunicorn)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/KeukaSensorProd/keuka
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/KeukaSensorProd/keuka
ExecStart=/home/pi/KeukaSensorProd/venv/bin/gunicorn -c /home/pi/KeukaSensorProd/keuka/gunicorn.conf.py wsgi:app
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target