# systemd/keuka-sensor.service
[Unit]
Description=Keuka Sensor (Gunicorn, STA+AP Support)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=pi
Group=pi
WorkingDirectory=/home/pi/KeukaSensorProd/keuka
# Load environment vars for STA/AP interface selection
EnvironmentFile=/etc/default/keuka-sensor
# Start Gunicorn with the new modular app
ExecStart=/home/pi/KeukaSensorProd/venv/bin/gunicorn \
    -c /home/pi/KeukaSensorProd/config/gunicorn.conf.py \
    'app:create_app()'
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
