# systemd/wifi-watchdog.service
[Unit]
Description=Keuka Wi-Fi Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
ExecStart=/bin/bash /home/pi/KeukaSensorProd/scripts/wifi-watchdog.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
